# 0CTF_TCTF-2021-Finals/babalogin/deployment/babalogin@.service
[Service]
ExecStart=-/usr/bin/timeout --foreground 1m /usr/bin/podman run --rm --init --security-opt seccomp=unconfined -i babalogin
StandardInput=socket
