# 0CTF_TCTF-2021-Finals/babalogin/deployment/babalogin.socket
[Socket]
ListenStream=32021
Accept=yes
MaxConnections=10000
MaxConnectionsPerSource=100

[Install]
WantedBy=sockets.target
