# 0CTF_TCTF-2021-Finals/babylogin/deployment/babylogin@.service
[Service]
ExecStart=-/usr/bin/timeout --foreground 2m /usr/bin/podman run --rm --init -i babylogin
StandardInput=socket
